# resource_monitor.service
[Unit]
Description=FreeLink Resource Monitor (CPU/RAM/Disk alerts)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/freelink
EnvironmentFile=/opt/freelink/.env
ExecStart=/opt/freelink/venv/bin/python3 /opt/freelink/resource_monitor.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
